# conf/supervisor/conf.d/plrfs.conf
[program:plrfs]
command = /opt/plrfs/groovy/rpcserver.sh
stdout_logfile = /tmp/supervisor-%(program_name)s.log
stderr_logfile = /tmp/supervisor-%(program_name)s.log
user = %(ENV_NB_USER)s
